# crons/iscsi_backup.cron
0 4 21 * * root export ELISE_ROOT_DIR='/mnt/tvault/es-labs/projects/elise'; ${ELISE_ROOT_DIR}/scripts/lab_backup.sh &>> /home/austin/backup_cron.log
