# crons/iscsi_backup.cron
0 4 15 * * austin /mnt/tvault/es-labs/projects/es-labs-bash/scripts/lab_backup.sh &>> /home/austin/backup_cron.log
